# package/base-files/files/etc/uci-defaults/99-remove-bad-opkg-sources
sed -ie '/immortalwrt_nas/d' /etc/opkg/distfeeds.conf 
sed -ie '/immortalwrt_nas_luci/d' /etc/opkg/distfeeds.conf 
sed -ie '/immortalwrt_istore/d' /etc/opkg/distfeeds.conf
sed -ie '/immortalwrt_istore_packages/d' /etc/opkg/distfeeds.conf 

exit 0