# package/base-files/files/etc/uci-defaults/99-remove-bad-opkg-sources
sed -ie '/immortalwrt_third/d' /etc/opkg/distfeeds.conf 

exit 0